# configuration/supervisor/conf.d/apache.conf
[program:apache]
priority=30
; directory=/tmp
command=/usr/sbin/apachectl -f /etc/apache2/apache2.conf
user=www-data
autostart=true
autorestart=true
