# configuration/supervisor/conf.d/nagios.conf
[program:nagios]
priority=20
; directory=/tmp
command=/usr/sbin/nagios3 -d /etc/nagios3/nagios.cfg
user=nagios
autostart=true
autorestart=true
