# configuration/supervisor/conf.d/nrpe.conf
[program:nrpe]
priority=40
; directory=/tmp
command=/usr/sbin/nrpe -c /etc/nagios/nrpe.cfg -d
user=nagios
autostart=true
autorestart=true
