# configuration/supervisor/conf.d/postfix.conf
[program:postfix]
priority=10
; directory=/tmp
command=/usr/sbin/postfix -c /etc/postfix
user=root
autostart=true
autorestart=true
